# system/kiosk.service
[Unit]
Description=PokelearnOS Kiosk Application
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pokelearnos
Group=pokelearnos
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pokelearnos/.Xauthority
Environment=DATABASE_URL=sqlite:///var/lib/pokelearnos/pokelearnos.db
Environment=APP_ENV=kiosk
Environment=HOST=127.0.0.1
Environment=PORT=8765
WorkingDirectory=/opt/pokelearnos
ExecStart=/opt/pokelearnos/system/kiosk-launcher.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pokelearnos

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/pokelearnos /var/lib/pokelearnos

[Install]
WantedBy=graphical.target
